# ops/systemd/561-group-site.service
[Unit]
Description=561.group public identity site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/emsenn/info/lib/emsenn/services/561-group/sites/561.group
Environment=SITE_561_GROUP_HOST=127.0.0.1
Environment=SITE_561_GROUP_PORT=18788
ExecStart=/usr/bin/env pnpm serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
